# systemd/openclaw.service
# systemd/openclaw.service — Production Service File
#
# Runs the ReClaw agent as a persistent systemd user service.
# Restarts automatically on failure. Logs go to journald.
#
# SETUP:
#   1. Copy this file to ~/.config/systemd/user/openclaw.service
#   2. Edit WorkingDirectory and ExecStart to match your actual paths
#   3. systemctl --user daemon-reload
#   4. systemctl --user enable openclaw
#   5. systemctl --user start openclaw
#   6. systemctl --user status openclaw
#
# VIEW LOGS:
#   journalctl --user -u openclaw -f

[Unit]
Description=ReClaw Personal AI Agent
After=network-online.target
Wants=network-online.target

[Service]
# --- Edit these two lines to match your setup ---
WorkingDirectory=/home/YOUR_USERNAME/reclaw-scaffold
ExecStart=/home/YOUR_USERNAME/reclaw-scaffold/venv/bin/python src/main.py
# ------------------------------------------------

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Load environment from .env file
EnvironmentFile=/home/YOUR_USERNAME/reclaw-scaffold/.env

# Keep service alive even if user logs out
# (requires: loginctl enable-linger YOUR_USERNAME)
KillMode=process

[Install]
WantedBy=default.target
