# systemd/reclaw-autoresearch.timer
# systemd/reclaw-autoresearch.timer — Overnight Optimization Loop Timer
#
# Triggers reclaw-autoresearch.service nightly at 2:00 AM local time.
# Runs the autoresearch orchestrator to evaluate and log agent quality.
#
# SETUP:
#   1. Copy this file to ~/.config/systemd/user/reclaw-autoresearch.timer
#   2. Copy reclaw-autoresearch.service to the same directory
#   3. systemctl --user daemon-reload
#   4. systemctl --user enable reclaw-autoresearch.timer
#   5. systemctl --user start reclaw-autoresearch.timer
#
# VERIFY IT'S SCHEDULED:
#   systemctl --user list-timers
#
# CHANGE THE TIME:
#   Edit OnCalendar below. Examples:
#     daily        = every day at midnight
#     *-*-* 02:00  = every day at 2am
#     Mon 03:00    = every Monday at 3am

[Unit]
Description=ReClaw Autoresearch Timer — runs nightly at 2am

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
